# dabplus-monitor.service
[Unit]
Description=DAB+ Monitor
After=network.target icecast2.service

[Service]
Environment=RTLSDR_PPM=-28
Type=simple
User=graffiti
WorkingDirectory=/home/graffiti/dabplus-monitor
ExecStartPre=+/bin/bash -c 'fuser -k 5000/tcp 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'pkill -9 welle-cli || true'
ExecStartPre=/bin/sleep 2
ExecStart=/home/graffiti/dabplus-monitor/venv/bin/gunicorn \
  --workers 1 --threads 16 --worker-class gthread \
  --bind 0.0.0.0:5000 --timeout 300 \
  --error-logfile /home/graffiti/dabplus-monitor/logs/gunicorn-error.log \
  --capture-output \
  --log-level info \
  app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
